# etc/supd/apps/tpl/lotus-worker-t28-2.ini
[program:lotus-worker-t28-2]
environment=GIN_MODE="release"
directory=$PRJ_ROOT/apps/lotus
command=./worker-t28-2.sh /data/sdb/lotus-user-1/.lotusminer /data/cache/.lotusworker
autostart=false
startsecs=3
startretries=3
autorestart=true
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=$PRJ_ROOT/var/log/lotus-worker-t28-2.logfile.stdout
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=10
stderr_logfile=$PRJ_ROOT/var/log/lotus-worker-t28-2.logfile.stderr
stderr_logfile_maxbytes=500MB
stderr_logfile_backups=10

