# etc/supd/apps/tpl/lotus-storage-2.ini
[program:lotus-storage-2]
environment=GIN_MODE="release"
directory=$PRJ_ROOT/apps/lotus
command=./storage.sh $PRJ_ROOT/var/lotus-storage-2/ /data/zfs2 :1350 :1351 :1352
autostart=false
startsecs=3
startretries=3
autorestart=true
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=$PRJ_ROOT/var/log/lotus-storage-2.logfile.stdout
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=10
stderr_logfile=$PRJ_ROOT/var/log/lotus-storage-2.logfile.stderr
stderr_logfile_maxbytes=500MB
stderr_logfile_backups=10

