# etc/supd/apps/tpl/lotus-datacap-chain.ini
[program:lotus-datacap-chain]
environment=GIN_MODE="release"
directory=$PRJ_ROOT/apps/lotus
command=./lotus-datacap-chain.sh /data/lotus-datacap/.lotus-proxy
autostart=false
startsecs=3
startretries=3
autorestart=true
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=$PRJ_ROOT/var/log/lotus-datacap-chain.logfile.stdout
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=10
stderr_logfile=$PRJ_ROOT/var/log/lotus-datacap-chain.logfile.stderr
stderr_logfile_maxbytes=500MB
stderr_logfile_backups=10

